# deploy/doge5m.service
[Unit]
Description=BTC5MIN Polymarket collector (DOGE Up/Down 5m)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/btc5min
ExecStart=/opt/btc5min/venv/bin/python3 -m src.main --asset doge --interval 5m
Restart=always
RestartSec=5
StandardOutput=append:/opt/btc5min/data/collector-doge5m.log
StandardError=journal

[Install]
WantedBy=multi-user.target
